# docker/etc/services.d/dbus/run
#!/usr/bin/with-contenv sh

exec dbus-daemon --system --nofork